# 02-stage-audiodriver-respeaker_lite/01-keep-audio-alive/files/keep-audio-alive.service
[Unit]
Description=Keep soundcard alive (play silence)
After=user@1000.service sound.target alsa-restore.service user@1000.service:pipewire.service 

[Service]
Type=simple
ExecStart=/usr/bin/keep-audio-alive.sh
Restart=always
User=root
StandardOutput=null

[Install]
WantedBy=multi-user.target
